# football/football-live.service
[Unit]
Description=Football Live Match Monitor
After=network.target

[Service]
ExecStart=/usr/bin/python3 "/root/CascadeProjects/sports bot/football/live.py" --continuous --interval 30
WorkingDirectory="/root/CascadeProjects/sports bot"
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
User=root
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
